# preflight/services/antimicrox-hyprdaemon.service
[Unit]
Description=Antimicrox Hyprland Daemon

[Service]
Type=simple
ExecStart=/usr/local/bin/antimicrox-hyprdaemon
Restart=on-failure
RestartSec=5
